# include/etc/cron.d/fetchmail
#
#  SpamAssassin Bayes learning from mailboxes 
#
# m h dom mon dow user  command
#1,31 *	* * *	root	if [ "0$USE_FETCHMAIL_PL" != "0" ]; then /bin/fetchmail.pl 2>&1; fi
1,31 *	* * *	root	if [ "0$USE_FETCHMAIL_PL" != "0" ]; then echo /bin/fetchmail.pl 2>&1; fi

