# include/etc/cron.d/sa-learn
#
#  SpamAssassin Bayes learning from mailboxes 
#
# m h dom mon dow user  command
30 0,3,6,9,12,15,18,21	* * *	root	/bin/spam-learn.sh 2>&1

